# roles/service/files/hello.service
[Unit]
Description=hello service

[Service]
ExecStart=/usr/local/bin/hello.sh

[Install]
WantedBy=multi-user.target
